# docker/conf/supervisor.ini
[supervisord]
nodaemon=true
pidfile=/tmp/fttt_docker_backend_supervisord.pid

[program:web]
command=gunicorn fortytwo.wsgi -b 0.0.0.0:8000
directory=/app
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/webout.log
stderr_logfile=/var/log/weberr.log
